# debian/taler-exchange.taler-exchange.service
[Unit]
Description=GNU Taler exchange
After=postgres.service network.target

[Install]
WantedBy=multi-user.target
