# debian/taler-exchange.taler-exchange.slice
[Unit]
Description=Slice for GNU taler exchange processes
Before=slices.target

[Slice]
# Add settings that should affect all GNU Taler exchange
# components here.
